# debian/mastodon-image-poster.timer
[Unit]
Description=Post an image to Mastodon daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
