# debian/mastodon-image-poster/DEBIAN/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

case "$1" in
    configure)
        # Retrieve debconf values
        db_get mastodon-image-poster/image_directory
        IMAGE_DIR="$RET"
        db_get mastodon-image-poster/sort_order
        SORT_ORDER="$RET"
        db_get mastodon-image-poster/instance_url
        INSTANCE_URL="$RET"
        db_get mastodon-image-poster/access_token
        ACCESS_TOKEN="$RET"
        db_get mastodon-image-poster/status_text
        STATUS_TEXT="$RET"

        # Create config directory and file
        mkdir -p /etc/mastodon-image-poster
        cat > /etc/mastodon-image-poster/config.ini <<EOF
[mastodon]
instance_url = ${INSTANCE_URL}
access_token = ${ACCESS_TOKEN}

[images]
directory = ${IMAGE_DIR}
sort_order = ${SORT_ORDER}
status_text = ${STATUS_TEXT}
EOF
        chmod 600 /etc/mastodon-image-poster/config.ini

        # Create state directory
        mkdir -p /var/lib/mastodon-image-poster

        # Make the script executable
        chmod +x /usr/bin/mastodon_image_poster.py

        # Enable and start the timer
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable mastodon-image-poster.timer
            systemctl start mastodon-image-poster.timer
        fi
        ;;
esac

# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'mastodon-image-poster.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'mastodon-image-poster.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'mastodon-image-poster.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'mastodon-image-poster.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'mastodon-image-poster.timer' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'mastodon-image-poster.timer'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'mastodon-image-poster.timer' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'mastodon-image-poster.timer' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'mastodon-image-poster.service' 'mastodon-image-poster.timer' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
